# configs/feather.service
[Unit]
Description=Feather's GRPC Unix Socket Server
After=network.target

[Service]
ExecStart=
User=feather
Group=feather